# rootfs/etc/supervisor.d/crond.ini
[program:crond]
command=/usr/sbin/crond -m %(ENV_APP_USER)s -L 15 -f
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisord/crond.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
