# rootfs/etc/supervisor.d/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -Dep 22
numprocs=1
autostart=true
autorestart=true
